# postfix/3.5/config/supervisor.d/rsyslog.ini
[program:rsyslogd]
command = /usr/sbin/rsyslogd -n
autorestart = false
startsecs = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
