# postfix/3.2/root/etc/supervisor.d/postfix.ini
[program:postfix]
command = /usr/sbin/postfix -c /etc/postfix start
autorestart = false
directory = /etc/postfix
process_name = master
startsecs = 0
